# supervisord.conf
[program:plone]
command = /home/dotcloud/current/bin/plone console
process_name = plone
directory = /home/dotcloud/current/bin
priority = 20
redirect_stderr = true
user = dotcloud
